# socket_activation/auto_active@.socket
[Unit]
Description=auto socket %i

[Socket]
ListenStream=99%i
BindIPv6Only=both

[Install]
WantedBy=sockets.target
